# srcs/influxdb/supervisord.conf
[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
loglevel=debug                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon

[program:influxdb]
command= influxd
autorestart=true
autostart=true

