# tmp/phpmyadmin/supervisord.conf
[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
loglevel=debug                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon

[program:nginx]
command= /usr/sbin/nginx -g 'daemon off;' 
autorestart=true
autostart=true

[program:php7]
command= /usr/sbin/php-fpm7 -D
autorestart=true
autostart=true

